# autoscreenshot.path
[Path]
PathModified=%h/.local/share/Steam/userdata/ACCOUNT_ID/760/screenshots.vdf
Unit=autoscreenshot.service
